# NBA_Schedules/crontab
5 00 1 */12 * python3 /bs4_scraper.py >> /logfile

